# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:icecast]
command=/usr/bin/icecast2 -c /etc/icecast2/icecast.xml
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/icecast.log
stderr_logfile=/var/log/supervisor/icecast_err.log
priority=10

[program:liquidsoap]
command=/usr/bin/liquidsoap /app/config/liquidsoap.liq
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/liquidsoap.log
stderr_logfile=/var/log/supervisor/liquidsoap_err.log
priority=20
startsecs=5

[program:flask]
command=gunicorn --worker-class eventlet -w 1 --bind 127.0.0.1:5000 "app:create_app()"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/flask.log
stderr_logfile=/var/log/supervisor/flask_err.log
priority=30

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
priority=40
